# files-portal/droplet-sync/kb-sync.timer
[Unit]
Description=Run KB sync every 60 seconds

[Timer]
OnBootSec=60s
OnUnitActiveSec=60s
Unit=kb-sync.service

[Install]
WantedBy=timers.target
